# docker/supervisor/etc/supervisor/conf.d/messenger-worker.conf
[program:messenger-email_notification]
command=php /var/www/html/app/bin/console messenger:consume email_notification --time-limit=3600
user=root
numprocs=2
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d

[program:messenger-loggable]
command=php /var/www/html/app/bin/console messenger:consume loggable --time-limit=3600
user=root
numprocs=2
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d
